# noblenest-academy/.docker/supervisord.conf
[supervisord]
nodaemon = true
user = root
logfile = /var/log/supervisor/supervisord.log
pidfile = /var/run/supervisord.pid

; PHP-FPM
[program:php-fpm]
process_name = %(program_name)s
command = /usr/local/sbin/php-fpm -F
autostart = true
autorestart = true
stdout_logfile = /var/log/supervisor/php-fpm.log
stderr_logfile = /var/log/supervisor/php-fpm.log
user = www-data

; Nginx
[program:nginx]
process_name = %(program_name)s
command = /usr/sbin/nginx -g "daemon off;"
autostart = true
autorestart = true
stdout_logfile = /var/log/supervisor/nginx.log
stderr_logfile = /var/log/supervisor/nginx.log
user = root

; Laravel Horizon (async job processor)
[program:horizon]
process_name = %(program_name)s
command = php /app/artisan horizon
directory = /app
autostart = true
autorestart = true
user = www-data
redirect_stderr = true
stdout_logfile = /var/log/supervisor/horizon.log
numprocs = 1
priority = 999

; Event cache warmer (optional — runs every hour)
[program:event-cache]
process_name = %(program_name)s
command = php /app/artisan event:cache
directory = /app
autostart = false
autorestart = false
startsecs = 0
user = www-data
redirect_stderr = true
stdout_logfile = /var/log/supervisor/event-cache.log

[group:laravel-workers]
programs = horizon
priority = 999
