# noblenest-academy/scripts/supervisor/noblenest-curriculum-ai.conf
[program:noblenest-curriculum-ai]
process_name=%(program_name)s
command=uvicorn main:app --host 127.0.0.1 --port 8001 --workers 2
directory=/var/www/noblenest/services/curriculum-ai
autostart=true
autorestart=true
user=deploy
redirect_stderr=true
stdout_logfile=/var/www/noblenest/storage/logs/curriculum-ai.log
